# services/camdisplay.timer
[Unit]
Description=Timer for the camdisplay service

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
